# grocery.service
[Unit]
Description=Grocery List Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/dino/grocery
ExecStart=/home/dino/grocery/venv/bin/python /home/dino/grocery/bot.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=grocery

[Install]
WantedBy=default.target
